# templates/vykar.container.j2
[Unit]
Description=Vykar Backup Server
After=network-online.target
RequiresMountsFor={{ services.vykar.zfs_path }}

[Container]
Image={{ services.vykar.image }}
ContainerName=vykar
# Bridge networking for UFW compatibility
Network=podman
PublishPort={{ services.vykar.container_port }}:8080

# Environment
Environment=VYKAR_HOSTNAME={{ services.vykar.domain }}
Environment=USER_UID=1000
Environment=USER_GID=1000

# Storage: Config on SSD, Backups on 20TB ZFS
Volume={{ services.vykar.ssd_path }}:/config
Volume={{ services.vykar.zfs_path }}:/data

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
